# root/etc/services.d/qbittorrent/run
#!/bin/sh
exec s6-setuidgid bsd /usr/local/bin/qbittorrent-nox \
    --confirm-legal-notice \
    --profile=/config
